# scripts/supervisord-socat.conf
[supervisord]
nodaemon=true

[program:socat]
command=/start_socat.sh

[program:registrator]
command=/usr/src/app/_build/default/rel/registrator/bin/registrator foreground